# etc/btrfs-auto-snapshot.cron.hourly
#!/bin/sh
PATH="$PATH:/usr/local/sbin"

# Only call btrfs-auto-snapshot if it's available
which btrfs-auto-snapshot > /dev/null || exit 0

exec btrfs-auto-snapshot --quiet --syslog --label=hourly --keep=24 //
